# examples/claude-code/com.user.rapture-processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.rapture-processor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "$HOME/Documents/Rapture Notes" &amp;&amp; ls *.txt &gt;/dev/null 2&gt;&amp;1 &amp;&amp; /usr/local/bin/claude --model haiku -p "process new notes per ./CLAUDE.md"</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/Users/YOUR_USERNAME/Documents/Rapture Notes</string>
  </array>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/rapture-processor.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/rapture-processor.err.log</string>
</dict>
</plist>
